# tts-websocket.service
[Unit]
Description=TTS PHP Ratchet WebSocket Server (wss://mentalfitness.store:8091)
After=network.target mysql.service
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
ExecStartPre=-/bin/bash -c 'fuser -k 8091/tcp 2>/dev/null || true'
ExecStart=/usr/bin/php /var/www/clients/client1/web71/web/artisan tts:websocket --port=8091
WorkingDirectory=/var/www/clients/client1/web71/web
Restart=always
RestartSec=10
User=web71
Group=client1
Environment=PATH=/usr/bin:/usr/local/bin
Environment=APP_ENV=production

[Install]
WantedBy=multi-user.target
